# sources/meta-bober/recipes-core/add-data-part/files/add-data-part.service
[Unit]
Description=Add data partition

[Service]
Type=oneshot
ExecStart=/usr/bin/add-data-part start

[Install]
WantedBy=multi-user.target