# sources/meta-bober/recipes-core/add-data-part/files/resize-datafs.service
[Unit]
Description=Resize data partition

[Service]
Type=oneshot
ExecStart=/usr/bin/resize-datafs start

[Install]
WantedBy=multi-user.target