# sources/meta-bober/recipes-app/app/files/qt-app.service
[Unit]
Description=Qt Test App

[Service]
ExecStart=/usr/bin/cube -platform eglfs

[Install]
WantedBy=sysinit.target