# sources/meta-bober/recipes-app/app/files/qt-app.service
[Unit]
Description=Qt Test App

[Service]
ExecStart=/usr/bin/qt_test -platform eglfs

[Install]
WantedBy=multi-user.target